# agent/backups.cron

################# AGENT CRON FILE ###################

# Do a incremental backup every minute
* * * * * root SERVER=$SERVER /copy_all_volumes_except_untagged.sh

