# agent/backups.cron

################# AGENT CRON FILE ###################

# Do a incremental backup every minute
* * * * * root /usr/bin/rsync -a --delete /volumes disasterman:/backups/latest/`hostname`

